# server/caddy.service
[Unit]
Description=Toolbox.io website caddy
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/caddy run
WorkingDirectory=/root/site/server

[Install]
WantedBy=multi-user.target